# deploy/instances/rel-ladder/rel-ladder.service
[Unit]
Description=Bughouse Chess Ladder - Release
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/bughouse-ladder/deploy/instances/rel-ladder
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/var/www/bughouse-ladder/server/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
